# systemd/rabbitcrm-backend.service
[Unit]
Description=RabbitCRM Backend (Django)
After=network.target

[Service]
Type=simple
User=himik
WorkingDirectory=/home/himik/rabbit-crm/backend
Environment="PATH=/home/himik/rabbit-crm/backend/.venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/himik/rabbit-crm/backend/.venv/bin/python /home/himik/rabbit-crm/backend/manage.py runserver 0.0.0.0:8000
Restart=on-failure
RestartSec=5
StandardOutput=append:/var/log/rabbitcrm-backend.log
StandardError=append:/var/log/rabbitcrm-backend.log

[Install]
WantedBy=multi-user.target
